# supervisord.conf
[supervisord]
nodaemon=true

[program:listener]
command=python app/listener.py
autorestart=true
stdout_logfile=app/log.txt
stderr_logfile=app/log.txt
